# package/trivialdns_init.conf
# vim: set ft=upstart ts=4 et:
description "trivialdns"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5

exec /usr/local/bin/trivialdns
